# target/audiovisio.ini
[program:audiovisio-api]
command=/usr/bin/java -jar /var/www/virtual/guido/api.mixing-senses.art/vi.son-audiovisio-db-0.0.1.jar
startsecs=60